# usr/local/opendomo/daemons/opencvodos.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides:          opencvodos
# Required-Start:    
# Required-Stop:
# Should-Start:      
# Default-Start:     1 2 3 4 5
# Default-Stop:      0 6
# Short-Description: Vision
# Description:       Vision
#
### END INIT INFO

# Copyright(c) 2014 OpenDomo Services SL. Licensed under GPL v3 or later
